# alpine/packages/containerd/etc/init.d/containerd
#!/sbin/openrc-run

depend()
{
	after docker
}

start()
{
	ebegin "Running system containerd"

	# set ulimits as high as possible
	ulimit -n 1048576
	ulimit -p unlimited

	[ ! -e /usr/bin/containerd ] && \
	( cd /usr/bin && \
	  ln -s docker-runc runc && \
	  ln -s docker-containerd-shim containerd-shim && \
	  ln -s docker-containerd-ctr containerd-ctr && \
	  ln -s docker-containerd containerd \
	)

	/usr/bin/containerd 1>&2 2>/var/log/containerd.log &

	ewaitfile 5 /var/run/containerd/containerd.sock

	eend $? "Failed to start system containerd"
}
