# alpine/packages/diagnostics/etc/init.d/diagnostics
#!/sbin/openrc-run

depend()
{
	after docker containerd
}

start()
{
	/usr/bin/diagnostics
}
