# alpine/packages/diagnostics/etc/init.d/diagnostics
#!/sbin/openrc-run

depend()
{
	after docker containerd
}

start()
{
	ebegin "Checking system state"

	DIAGNOSTICS_SERVER_FLAGS=""

	case "$(mobyplatform)" in
		"aws")
			DIAGNOSTICS_SERVER_FLAGS="-http"
		;;
		"azure")
			DIAGNOSTICS_SERVER_FLAGS="-http"
		;;
		"windows")
			DIAGNOSTICS_SERVER_FLAGS="-hvsock"
		;;
		"mac")
			DIAGNOSTICS_SERVER_FLAGS="-vsock"
		;;
	esac

	/usr/bin/diagnostics-server ${DIAGNOSTICS_SERVER_FLAGS} &
	/usr/bin/diagnostics
}
