# alpine/packages/diagnostics/etc/init.d/diagnostics
#!/sbin/openrc-run

depend()
{
	after docker mdnstool hupper
}

start()
{
	ebegin "Checking system state"

	/usr/bin/diagnostics-server &
	/usr/bin/diagnostics
}
