# alpine/packages/diagnostics/etc/init.d/diagnostics-server
#!/sbin/openrc-run

start()
{
	ebegin "Starting diagnostics server"

	DIAGNOSTICS_SERVER_FLAGS=""

	case "$(mobyplatform)" in
		"aws")
			DIAGNOSTICS_SERVER_FLAGS="-http"
		;;
		"azure")
			DIAGNOSTICS_SERVER_FLAGS="-http"
		;;
		"gcp")
			DIAGNOSTICS_SERVER_FLAGS="-http"
		;;
		"windows")
			DIAGNOSTICS_SERVER_FLAGS="-hvsock"
		;;
		"mac")
			DIAGNOSTICS_SERVER_FLAGS="-vsock"
		;;
	esac

	/usr/bin/diagnostics-server ${DIAGNOSTICS_SERVER_FLAGS} &

	eend 0
}
