# alpine/packages/hupper/etc/init.d/hupper
#!/sbin/openrc-run

description="docker config update manager"

depend()
{
	after docker
}

start()
{
	ebegin "Starting docker config update manager"

	PIDFILE=/run/hupper.pid
	DOCKERPIDFILE=/run/docker.pid
	WATCH_CONFIG=$(mobyconfig watch /etc/docker/daemon.json)
	WATCH_PROXY=$(mobyconfig watch proxy)

	[ -z "${WATCH_CONFIG}" ] && exit 1
	[ -z "${WATCH_PROXY}" ] && exit 1

	start-stop-daemon --start --quiet \
		--background \
		--exec /bin/hupper \
		--make-pidfile --pidfile ${PIDFILE} \
		-- -pidfile ${PIDFILE} -huppidfile ${DOCKERPIDFILE} \
		   -path ${WATCH_CONFIG} -path ${WATCH_PROXY}

	eend $? "Failed to start hupper"
}

stop()
{
	ebegin "Stopping docker config update manager"

	PIDFILE=/run/hupper.pid

	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}

	eend $? "Failed to stop hupper"
}
