# alpine/packages/proxy/etc/init.d/proxy
#!/sbin/openrc-run

depend()
{
        before docker
}

start()
{
	ebegin "Setting up proxy port mount"

	mkdir -p /port
	mount -t 9p -o trans=virtio,dfltuid=1001,dfltgid=50,version=9p2000 port /port

	eend $? "Failed to mount proxy port filesystem"
}
