# alpine/etc/inittab
# /etc/inittab

::sysinit:/sbin/rc sysinit
::wait:/sbin/rc default

# Put a getty on the serial port
ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/sbin/rc shutdown
