# alpine/etc/periodic/15min/trim
#!/bin/sh

# Windows and Mac have TRIM support, clean out unused frequently
case "$(mobyplatform)" in
  "windows"|"mac")
     /sbin/fstrim /var
     ;;
esac
