# alpine/packages/transfused/etc/init.d/transfused
#!/sbin/openrc-run

description="fuse proxy server"

depend()
{
	need 9pinit
}

start()
{
	[ -d /Transfuse ] || exit 0

	ebegin "Starting FUSE socket passthrough"

        if cat /proc/cmdline | grep -q 'com.docker.driverDir'
        then
                DRIVERDIR="/Mac$(cat /proc/cmdline | sed -e 's/.*com.docker.driverDir="//' -e 's/".*//')"
                RUNTIME_LOGFILE=${DRIVERDIR}/log/transfused.log
        else
                ID=$(LC_CTYPE=C tr -dc A-Za-z0-9 < /dev/urandom | head -c 8)
                RUNTIME_LOGFILE="/Mac/private/tmp/transfused_${ID}.log"
        fi

        PIDFILE=/var/run/transfused.pid
        STARTUP_LOGFILE=/var/transfused_start.log
        MOUNT_TRIGGER=/Mac

	start-stop-daemon --start --quiet \
		--background \
		--exec /sbin/transfused \
		--pidfile ${PIDFILE} \
		-- \
                -p "${PIDFILE}" \
                -l "${STARTUP_LOGFILE}" \
                -m "${MOUNT_TRIGGER}" \
                -t "${RUNTIME_LOGFILE}"

	eend $? "Failed to start transfused"
}

stop()
{
	[ -d /Transfuse ] || exit 0

	ebegin "Stopping FUSE socket passthrough"

	PIDFILE=/var/run/transfused.pid

	start-stop-daemon --stop --quiet \
		--pidfile "${PIDFILE}"

	eend $? "Failed to stop transfused"
}
