# alpine/packages/transfused/etc/init.d/transfused
#!/sbin/openrc-run

description="fuse proxy server"

depend()
{
	need 9pinit
}

start()
{
	[ -d /Transfuse ] || exit 0

	ebegin "Starting FUSE socket passthrough"

	PIDFILE=/var/run/transfused.pid
	LOGFILE=/vatr/log/transfused.log

	start-stop-daemon --start --quiet \
		--background \
		--exec /sbin/transfused \
		--pidfile ${PIDFILE} \
		--stderr "${LOGFILE}" \
		--stdout "${LOGFILE}" \
		-- -p "${PIDFILE}"

	eend $? "Failed to start transfused"
}

stop()
{
	[ -d /Transfuse ] || exit 0

	ebegin "Stopping FUSE socket passthrough"

	PIDFILE=/var/run/transfused.pid

	start-stop-daemon --stop --quiet \
		--pidfile "${PIDFILE}"

	eend $? "Failed to stop transfused"
}
