# alpine/etc/init.d/fsck
#!/sbin/openrc-run
# do nothing as we do this in automount script

start()
{
	return 0
}

stop()
{
	return 0
}
