# alpine/packages/test/etc/init.d/test
#!/sbin/openrc-run

depend()
{
	need docker containerd
}

start()
{
	[ -d /test ] || exit 0

	ebegin "Running tests"

	if ctr run --bundle /test test
	then
		printf "Moby test suite PASSED\n"
	else
		printf "Moby test suite FAILED\n"
	fi

	eend "Tests completed"

	# now terminate with extreme prejudice
	poweroff -f
}
