# alpine/packages/hostsettings/etc/init.d/hostsettings
#!/sbin/openrc-run

description="Configuring settings from database."

depend() {
	before sysctl net
}

start() {
	ebegin "Configuring host settings from database"

	mobyconfig exists etc/sysctl.conf && mobyconfig get etc/sysctl.conf > /etc/sysctl.conf
	mobyconfig exists etc/sysfs.conf && mobyconfig get etc/sysfs.conf > /etc/sysfs.conf
	mobyconfig exists etc/resolv.conf && mobyconfig get etc/resolv.conf > /etc/resolv.conf

	mobyconfig exists etc/hosts && mobyconfig get etc/hosts >> /etc/hosts

	mobyconfig exists etc/ssl/certs/ca-certificates.crt && mobyconfig get etc/ssl/certs/ca-certificates.crt >> /etc/ssl/certs/ca-certificates.crt 

	eend 0
}
