# alpine/packages/hostsettings/etc/init.d/hostsettings
#!/sbin/openrc-run

# shellcheck disable=SC2034
description="Configuring settings from database."

depend() {
	before sysctl net
}

start() {
	ebegin "Configuring host settings from database"

	mobyconfig exists etc/sysctl.conf && mobyconfig get etc/sysctl.conf > /etc/sysctl.conf
	mobyconfig exists etc/sysfs.conf && mobyconfig get etc/sysfs.conf > /etc/sysfs.conf
	mobyconfig exists etc/resolv.conf && mobyconfig get etc/resolv.conf > /etc/resolv.conf
	mobyconfig exists random-seed && mobyconfig get random-seed > /dev/urandom

	mobyconfig exists etc/hosts && mobyconfig get etc/hosts >> /etc/hosts

	# handle static network config if configured
	mobyconfig exists net/config && NETCONFIG=$(mobyconfig get net/config)
	if [ "${NETCONFIG}" = "static" ]; then
	     # assume that the other configsDB entries exist
	     IP=$(mobyconfig get net/address)
	     MASK=$(mobyconfig get net/netmask)
	     GW=$(mobyconfig get net/gateway)

	     cp /etc/network/interfaces.template /etc/network/interfaces
	     {
                 echo "auto eth0"
                 echo "iface eth0 inet static"
                 echo "    address ${IP}"
                 echo "    netmask ${MASK}"
                 echo "    gateway ${GW}"
                 echo "    metric  200"
             } >> /etc/network/interfaces
	fi
	
	mobyconfig exists etc/ssl/certs/ca-certificates.crt && mobyconfig get etc/ssl/certs/ca-certificates.crt >> /etc/ssl/certs/ca-certificates.crt 

	eend 0
}
