# alpine/packages/hostsettings/etc/init.d/hostsettings
#!/sbin/openrc-run

description="Configuring settings from database."

depend() {
	before sysctl
}

start() {
	ebegin "Configuring host settings from database"

	mobyconfig exists etc/sysctl.conf && mobyconfig get etc/sysctl.conf > /etc/sysctl.conf

	# shift logs onto host before docker starts
	# busybox reopens its log files every second
	if cat /proc/cmdline | grep -q 'com.docker.driverDir'
	then
		mount --bind "${DRIVERDIR}/log" /var/log
	fi

        eend 0
}
