# alpine/packages/hostsettings/etc/init.d/hostsettings
#!/sbin/openrc-run

description="Configuring settings from database."

depend() {
	need transfused
	before sysctl net
}

start() {
	ebegin "Configuring host settings from database"

	mobyconfig exists etc/sysctl.conf && mobyconfig get etc/sysctl.conf > /etc/sysctl.conf

	# shift logs onto host before docker starts
	# busybox reopens its log files every second
	if cat /proc/cmdline | grep -q 'com.docker.driverDir'
	then
		DRIVERDIR="/Mac$(cat /proc/cmdline | sed -e 's/.*com.docker.driverDir="//' -e 's/".*//')"
		mount --bind "${DRIVERDIR}/log" /var/log
	fi

	mobyconfig exists etc/resolv.conf && mobyconfig get etc/resolv.conf > /etc/resolv.conf

	# handle static network config if configured
	mobyconfig exists net/config && NETCONFIG=`mobyconfig get net/config`
	if [ "${NETCONFIG}" = "static" ]; then
	     # assume that the other configsDB entries exist
	     IP=`mobyconfig get net/address`
	     MASK=`mobyconfig get net/netmask`
	     GW=`mobyconfig get net/gateway`

	     cp /etc/network/interfaces.template /etc/network/interfaces
	     echo                          >> /etc/network/interfaces
	     echo "auto eth0"              >> /etc/network/interfaces
	     echo "iface eth0 inet static" >> /etc/network/interfaces
	     echo "    address ${IP}"      >> /etc/network/interfaces
	     echo "    netmask ${MASK}"    >> /etc/network/interfaces
	     echo "    gateway ${GW}"      >> /etc/network/interfaces
	     echo "    metric  200"        >> /etc/network/interfaces
	fi

	eend 0
}
