# alpine/packages/hostsettings/etc/init.d/hostsettings
#!/sbin/openrc-run

description="Configuring settings from database."

depend() {
	# AWS gets settings via network
	need net
	before docker
}

start() {
	ebegin "Configuring host settings from database"

	mobyconfig exists random-seed && mobyconfig get random-seed > /dev/urandom

	mobyconfig exists etc/resolv.conf && mobyconfig get etc/resolv.conf > /etc/resolv.conf
	mobyconfig exists etc/hosts && mobyconfig get etc/hosts >> /etc/hosts

	mobyconfig exists etc/ssl/certs/ca-certificates.crt && mobyconfig get etc/ssl/certs/ca-certificates.crt >> /etc/ssl/certs/ca-certificates.crt 

	eend 0
}
