# alpine/packages/hostsettings/etc/init.d/hostsettings
#!/sbin/openrc-run

description="Configuring settings from database."

depend() {
	before hostname sysctl
}

start() {
	ebegin "Configuring host settings from database"

	mobyconfig exists etc/hostname && echo $(mobyconfig get etc/hostname) > /etc/hostname
	mobyconfig exists etc/sysctl.conf && mobyconfig get etc/sysctl.conf > /etc/sysctl.conf

        eend 0
}
