# alpine/etc/periodic/weekly/trim
#!/bin/sh

# TODO check with hdparm if supports TRIM
/sbin/fstrim /var || true
