# alpine/packages/oom/etc/init.d/oom
#!/sbin/openrc-run

# shellcheck disable=SC2034
description="oom killer settings"

depend()
{
	after docker
}

start()
{
	ebegin "Adjusting oom killer settings"

  find /run -name '*.pid' > tmp
	while IFS= read -r f; do
		PID=$(cat "$f")
		[ -f /proc/"$PID"/oom_score_adj ] && echo "-800" > /proc/"$PID"/oom_score_adj	
	done < tmp
  rm tmp
	eend $? "Failed to adjust oom settings"
}
