# alpine/packages/aws/etc/init.d/aws
#!/sbin/openrc-run

description="Bootstrap procedure if running on Docker for AWS"

depend()
{
	need docker
}

start()
{
	[ "$(mobyplatform)" != "aws" ] && exit 0
	ebegin "Running AWS-specific initialization"

	# Needed for ELB integration.
	mkdir -p /var/lib/docker/swarm

	# Get user data file and use it to bootstrap Moby in the cloud
	mobyconfig get userdata > /tmp/user-data

	# For now we will have a shell script which executes on boot.
	# TODO(nathanleclaire/kencochrane): Migrate this to mobyconfig, or similar.
	if [ $? -eq 0 ]
	then
		chmod 755 /tmp/user-data
		/tmp/user-data || true
	fi

	eend 0
}
