# alpine/packages/docker/etc/init.d/docker
#!/sbin/openrc-run

export DOCKER_RAMDISK="true"

command="${DOCKER_BINARY:-/usr/bin/docker}"

pidfile="/run/docker.pid"

if cat /proc/cmdline | grep -q 'com.docker.database'
then
	DATABASE="$(cat /proc/cmdline | sed -e 's/.*com.docker.database="//' -e 's/".*//')"
	CONFIG_FILE="/Database/branch/master/ro/${DATABASE}/etc/docker/daemon.json"
	[ -s ${CONFIG_FILE} ] && DOCKER_OPTS="${DOCKER_OPTS} --config-file ${CONFIG_FILE}"
fi

command_args="daemon -p \"${pidfile}\" ${DOCKER_OPTS}"

DOCKER_LOGFILE="/var/log/docker.log}"

start_stop_daemon_args="--background \
        --stderr \"${DOCKER_LOGFILE}\" --stdout \"${DOCKER_LOGFILE}\""

depend()
{
        after 9pudc 9pudfuse
        before chronyd
}
