# alpine/packages/docker/etc/init.d/docker
#!/sbin/openrc-run

export DOCKER_RAMDISK="true"

# choose a graph driver based on what is available and a preference
# if the user has not set any opts
# we do not yet provide easy access to /etc/conf.d/docker but may do
# otherwise let Docker choose
if cat /proc/filesystems | grep -q '\taufs$'
then
	DOCKER_OPTS="${DOCKER_OPTS:- -s aufs}"
elif cat /proc/filesystems | grep -q '\toverlay$'
then
	DOCKER_OPTS="${DOCKER_OPTS:- -s overlay}"
fi

# temporary fix until -rc2
DOCKER_OPTS="-s overlay"

command="${DOCKER_BINARY:-/usr/bin/docker}"

pidfile="/run/docker.pid"

if cat /proc/cmdline | grep -q 'com.docker.groupDir'
then
	GROUPDIR="/Mac$(cat /proc/cmdline | sed -e 's/.*com.docker.groupDir="//' -e 's/".*//')"
	[ -d "${GROUPDIR}" ] && DOCKER_LOGFILE="${GROUPDIR}/docker.log"
fi

if cat /proc/cmdline | grep -q 'com.docker.database'
then
	DATABASE="$(cat /proc/cmdline | sed -e 's/.*com.docker.database="//' -e 's/".*//')"
	CONFIG_FILE="/Database/branch/master/ro/${DATABASE}/etc/docker/daemon.json"
	[ -s ${CONFIG_FILE} ] && DOCKER_OPTS="${DOCKER_OPTS} --config-file ${CONFIG_FILE}"
fi

command_args="daemon -p \"${pidfile}\" ${DOCKER_OPTS}"

DOCKER_LOGFILE="${DOCKER_LOGFILE:-/var/log/${RC_SVCNAME}.log}"
start_stop_daemon_args="--background \
        --stderr \"${DOCKER_LOGFILE}\" --stdout \"${DOCKER_LOGFILE}\""

depend()
{
        after 9pudc
        before chronyd
}
