# pkg/init/etc/inittab
# /etc/inittab

::sysinit:/bin/rc.init

# Stuff to do for the 3-finger salute
::ctrlaltdel:/bin/rc.shutdown reboot

# Stuff to do on shutdown
::shutdown:/bin/rc.shutdown poweroff
