# boot2docker/rootfs/crontab
# resync with the host clock every hour
0 * * * * hwclock -s
