# alpine/packages/mdnstool/etc/init.d/mdnstool
#!/sbin/openrc-run

description="mDNS server"

depend()
{
        need 9pinit net
	after firewall
}

start()
{
	ebegin "Starting mDNS server"

	[ -n "${PIDFILE}" ] || PIDFILE=/var/run/mdnstool.pid

	export HOSTNAME=$(hostname -s).local.

	start-stop-daemon --start --quiet \
		--background \
		--exec /sbin/mdnstool \
		--make-pidfile --pidfile ${PIDFILE} \
		-- -if eth0 -hostname ${HOSTNAME}

	eend $? "Failed to start mDNS server"
}

stop()
{
	ebegin "Stopping mDNS server"

	[ -n "${PIDFILE}" ] || PIDFILE=/var/run/mdnstool.pid

	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}

	eend $? "Failed to stop mDNS server"
}
