# alpine/packages/mdnstool/etc/init.d/mdnstool
#!/sbin/openrc-run

description="mDNS server"

depend()
{
        need net
	after firewall
}

start()
{
        [ -d /sys/bus/vmbus ] && exit 0

        ebegin "Starting mDNS server"

	PIDFILE=/var/run/mdnstool.pid

	hostname -F /etc/hostname
	export HOSTNAME=$(hostname -s).local.

	start-stop-daemon --start --quiet \
		--background \
		--exec /sbin/mdnstool \
		--make-pidfile --pidfile ${PIDFILE} \
		-- -if eth0 -hostname ${HOSTNAME}

	eend $? "Failed to start mDNS server"
}

stop()
{
        [ -d /sys/bus/vmbus ] && exit 0

        ebegin "Stopping mDNS server"

	PIDFILE=/var/run/mdnstool.pid

	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}

	eend $? "Failed to stop mDNS server"
}
