# alpine/packages/mdnstool/etc/init.d/mdnstool
#!/sbin/openrc-run

description="mDNS server"

depend()
{
        need net
	after firewall
}

start()
{
        ebegin "Starting mDNS server"

	PIDFILE=/var/run/mdnstool.pid

	hostname -F /etc/hostname
	export HOSTNAME=$(hostname -s).local.

	INTF=eth0
	if [ "$(mobyconfig get network)" = "hybrid" ]; then
		INTF=eth1
	fi

	start-stop-daemon --start --quiet \
		--background \
		--exec /sbin/mdnstool \
		--make-pidfile --pidfile ${PIDFILE} \
		-- -if ${INTF} -hostname ${HOSTNAME}

	eend $? "Failed to start mDNS server"
}

stop()
{
        ebegin "Stopping mDNS server"

	PIDFILE=/var/run/mdnstool.pid

	start-stop-daemon --stop --quiet --pidfile ${PIDFILE}

	eend $? "Failed to stop mDNS server"
}
