# alpine/packages/hvtools/etc/init.d/hv_kvp_daemon
#!/sbin/openrc-run

HV_DAEMON=hv_kvp_daemon

depend()
{
        after dev
        need net
        before docker proxy vsudd diagnostics
}

start()
{
        [ "$(mobyplatform)" != "windows" ] && exit 0

        ebegin "Starting Hyper-V Daemon: ${HV_DAEMON}"

        # Delete the existing store
        rm -rf /var/lib/hyperv

        [ -n "${PIDFILE}" ] || PIDFILE=/var/run/${HV_DAEMON}.pid

        start-stop-daemon --start --quiet \
                --background \
                --exec /sbin/${HV_DAEMON} \
                --make-pidfile --pidfile ${PIDFILE} \
                --
        eend 0
}

stop()
{
        [ "$(mobyplatform)" != "windows" ] && exit 0

        ebegin "Stopping Hyper-V Daemon: ${HV_DAEMON}"

        [ -n "${PIDFILE}" ] || PIDFILE=/var/run/${HV_DAEMON}.pid

        start-stop-daemon --stop --quiet --pidfile ${PIDFILE}

        eend $? "Failed to stop ${HV_DAEMON}"
}
