# alpine/packages/azure/etc/init.d/azure
#!/sbin/openrc-run

description="Bootstrap procedure if running on Docker Azure edition"

depend()
{
	need docker
	need networking
}

initdns()
{
	while [ ! -f /var/lib/waagent/provisioned ]
	do
		sleep 5
	done

	# For use as a part of a sed expression to update hostname in
	# /etc/network/interfaces
	PATTERN='/iface eth0/a\
    hostname'

	# Hack: Drop 'hostname ${HOSTNAME}' into DHCP configuration file
	#
	# TODO(nathanleclaire): Indentation doesn't seem to work properly.  It
	# still works but how to enable this with shell script?
	sed -i "${PATTERN} $(hostname)" /etc/network/interfaces

	# Freshly configured to advertise the Azure-set hostname, obtain a new
	# DHCP lease.
	kill -SIGUSR1 $(pidof udhcp)
}

start()
{
	[ "$(mobyplatform)" != "azure" ] && exit 0
	ebegin "Running Azure-specific initialization"

	# TODO(nathanleclaire): With recent DHCP fixes is this still needed?
	for i in $(seq 1 20)
	do
		einfo "Pulling Windows Azure Linux Agent container"

		docker pull docker4x/agent-azure >/dev/null

		if [ $? -eq 0 ]
		then
			break
		fi

		# Wait for... network to come up?  DNS servers to be reachable?
		# Not certain, but Azure continually fails to achieve this pull so
		# far because it can't dial the DNS lookup properly.
		# 
		# TODO: Debug.
		sleep 5
	done

	einfo "Running Windows Azure Linux Agent container"

	docker run -d \
		--privileged \
		--name agent \
		--ipc host \
		--pid host \
		--net host \
		--restart unless-stopped \
		-v /usr/bin/docker:/usr/local/bin/docker:ro \
		-v /var/run/docker.sock:/var/run/docker.sock \
		-v /var/log:/var/log \
		-v /home:/home \
		-v /etc:/etc \
		-v /lib/modules:/lib/modules \
		-v /lib/firmware:/lib/firmware \
		-v /var/lib/waagent:/var/lib/waagent \
		docker4x/agent-azure

	# Wait for docker user to be added by agent.
	while [ ! -d /home/docker ]
	do
		sleep 5
	done

	# TODO: Make this cleaner.
	# User gets added by waagent.
	# Need to unlock it to login via SSH.
	passwd -u docker
	checkpath --directory --mode 0700 /home/docker/.ssh

	# Wait for custom data to arrive
	while [ ! -f /var/lib/waagent/CustomData ]
	do
		sleep 5
	done

	. /var/lib/waagent/CustomData

	initdns

	eend 0
}

stop()
{
	[ "$(mobyplatform)" != "azure" ] && exit 0
	docker rm -f agent || true
	passwd -l docker
}
