# alpine/packages/random/etc/init.d/random
#!/sbin/openrc-run

start()
{
	ebegin "Adding to random seed"

	/sbin/add-random

	eend $? "Failed to add randomness"
}
